# etc/systemd/system/backup-dvd-sr0.service
[Unit]
Description=Backup DVD

[Service]
Type=simple
ExecStart=/root/dvd-backup/dvd-backup.sh sr0
